# supervisord.d/laravel-queue.conf
[program:laravel-queue]
command=php /var/www/artisan queue:work --sleep=10 --tries=3 --daemon
autostart=true
autorestart=true
priority=20
